# olmanager-data-manager.service
# OLManager Data Manager — systemd service
# Place at: /etc/systemd/system/olmanager-data-manager.service
# Start:    sudo systemctl start olmanager-data-manager
# Enable:   sudo systemctl enable olmanager-data-manager
# Logs:     sudo journalctl -u olmanager-data-manager -f

[Unit]
Description=OLManager Data Manager (Next.js)
After=network.target

[Service]
Type=simple
User=debian
WorkingDirectory=/home/debian/olmanager-data-manager
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/node node_modules/.bin/next start
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=no
ReadWritePaths=/home/debian/olmanager-data-manager/src/data
ReadWritePaths=/tmp
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
